# temprecorder/debian/cron.d
* * * * * root /www/temprecorder/bin/record.sh
